# packaging/deb/DEBIAN/postinst
#!/bin/bash
# Polaris post-install: creates the polaris system user, sets up the
# capture root, installs indi-web into a dedicated venv, then enables
# and starts the systemd service.
#
# Runs after dpkg has unpacked /opt/polaris/, /lib/systemd/system/
# polaris.service, and /opt/polaris/appsettings.json. Idempotent so
# upgrades and re-installs do not break existing setups.

set -e

POLARIS_USER=polaris
POLARIS_GROUP=polaris
POLARIS_HOME=/home/polaris
POLARIS_FILES=$POLARIS_HOME/files
POLARIS_MODELS=$POLARIS_HOME/models
POLARIS_CONFIG_DIR=$POLARIS_HOME/.config/NINA.Polaris
INDIWEB_VENV=/opt/polaris-indiweb-venv

case "$1" in
    configure)
        # ---- 1. System user ----
        if ! getent passwd $POLARIS_USER >/dev/null; then
            adduser --system --group --home $POLARIS_HOME \
                    --shell /bin/bash --gecos "N.I.N.A. Polaris" $POLARIS_USER
            echo "polaris: created system user $POLARIS_USER"
        fi

        # adduser usually creates the home, but be defensive in case
        # someone removed it or pre-created the user differently.
        if [ ! -d $POLARIS_HOME ]; then
            mkdir -p $POLARIS_HOME
            chown $POLARIS_USER:$POLARIS_GROUP $POLARIS_HOME
            chmod 755 $POLARIS_HOME
        fi

        # ---- 2. Capture root ----
        if [ ! -d $POLARIS_FILES ]; then
            mkdir -p $POLARIS_FILES
            chown $POLARIS_USER:$POLARIS_GROUP $POLARIS_FILES
            echo "polaris: created capture root at $POLARIS_FILES"
        fi

        # ---- 2b. GraXpert ONNX models dir ----
        # OnnxModelRegistry checks /home/polaris/models as the Linux
        # default before falling back to the bundled wwwroot path.
        # Create empty so the operator (or a separate deploy script)
        # can rsync the GraXpert weights into it without sudo.
        if [ ! -d $POLARIS_MODELS ]; then
            mkdir -p $POLARIS_MODELS
            chown $POLARIS_USER:$POLARIS_GROUP $POLARIS_MODELS
            echo "polaris: created models dir at $POLARIS_MODELS"
        fi

        # ---- 3. Profile dir ----
        if [ ! -d $POLARIS_CONFIG_DIR ]; then
            mkdir -p $POLARIS_CONFIG_DIR
            chown -R $POLARIS_USER:$POLARIS_GROUP $POLARIS_HOME/.config
        fi

        # ---- 4. indi-web venv ----
        # One-time install. If user upgrades Polaris we leave the venv
        # alone; they can refresh indiweb manually with:
        #   sudo $INDIWEB_VENV/bin/pip install --upgrade indiweb
        if [ ! -x $INDIWEB_VENV/bin/indi-web ]; then
            echo "polaris: installing indi-web in $INDIWEB_VENV (one-time, ~30 seconds)"
            if ! python3 -m venv $INDIWEB_VENV; then
                echo "polaris: WARNING failed to create indi-web venv, embedded INDI manager will not work" >&2
                echo "polaris: install manually with: sudo python3 -m venv $INDIWEB_VENV && sudo $INDIWEB_VENV/bin/pip install indiweb" >&2
            else
                if ! $INDIWEB_VENV/bin/pip install --quiet --upgrade pip 2>/dev/null; then
                    echo "polaris: WARNING pip upgrade failed, continuing with bundled pip" >&2
                fi
                if ! $INDIWEB_VENV/bin/pip install --quiet indiweb; then
                    echo "polaris: WARNING indi-web pip install failed, embedded INDI manager will not work" >&2
                    echo "polaris: most likely no internet during install; retry with:" >&2
                    echo "polaris:   sudo $INDIWEB_VENV/bin/pip install indiweb legacy-cgi" >&2
                fi
                # Python 3.13 removed the 'cgi' stdlib module that
                # indi-web's vendored bottle.py imports unconditionally.
                # legacy-cgi is a PyPI backport that registers the
                # 'cgi' module name and makes the import succeed.
                # Harmless on Python 3.12 and earlier.
                if ! $INDIWEB_VENV/bin/pip install --quiet legacy-cgi; then
                    echo "polaris: WARNING legacy-cgi install failed; indi-web may fail to start on Python 3.13+" >&2
                fi
            fi
        fi

        # ---- 5. WiFi bootstrap script + PolicyKit ----
        # Mark the bootstrap script executable; dpkg-deb sometimes
        # loses the mode bit if the build host's umask is loose.
        # Idempotent: chmod is a no-op when the bits are already set.
        chmod 0755 /opt/polaris/bin/polaris-wifi-bootstrap.sh 2>/dev/null || true

        # PolicyKit reload picks up /etc/polkit-1/rules.d/50-polaris-nm.rules
        # which grants the polaris user nmcli access without sudo
        # prompts. polkit may not be installed on minimal containers;
        # the rule still lives on disk, ready for whenever polkit
        # arrives.
        if [ -d /run/systemd/system ]; then
            systemctl restart polkit.service >/dev/null 2>&1 || true
        fi

        # ---- 6. systemd ----
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload

            # Enable on first install; preserve state on upgrade
            if [ -z "$2" ]; then
                # $2 empty means fresh install (no previous version)
                systemctl enable polaris.service >/dev/null 2>&1 || true
                # First-boot WiFi bootstrap. Only enabled on a fresh
                # install so we do not stomp on a user who already
                # customised polaris-hotspot via the UI.
                systemctl enable polaris-wifi-bootstrap.service >/dev/null 2>&1 || true
            fi

            # Start (or restart if already running on upgrade)
            if systemctl is-active --quiet polaris.service; then
                systemctl restart polaris.service
            else
                systemctl start polaris.service
            fi
        else
            echo "polaris: systemd not active, skipping service enable/start"
            echo "polaris: when systemd is available run: sudo systemctl enable --now polaris.service"
        fi

        # ---- 7. Friendly summary ----
        HOST=$(hostname 2>/dev/null || echo "localhost")
        IP=$(hostname -I 2>/dev/null | awk '{print $1}' || echo "")
        echo ""
        echo "==============================================="
        echo " N.I.N.A. Polaris installed."
        echo ""
        echo " Open in your browser:"
        echo "   https://${HOST}.local:5000"
        if [ -n "$IP" ]; then
            echo "   https://${IP}:5000  (if mDNS does not resolve)"
        fi
        echo ""
        echo " First visit will show a 'self-signed certificate' warning"
        echo " (Polaris generates the cert at first launch). Accept it once"
        echo " per device and the browser remembers it. WebGPU acceleration"
        echo " for in-browser GraXpert AI inference needs HTTPS, hence this"
        echo " is the default."
        echo ""
        echo " Polaris also listens on http://127.0.0.1:5080 (loopback only,"
        echo " for SSH-tunneled access or a local reverse proxy fronting TLS"
        echo " elsewhere). LAN access is HTTPS:5000 only."
        echo ""
        echo " WiFi: on the next boot the Pi comes up as a hotspot named"
        echo "   SSID:     Polaris-Hotspot"
        echo "   Password: polaris1234"
        echo " Connect any device, open https://${HOST}.local:5000, then go to"
        echo " Settings → Network to switch the Pi onto your home WiFi instead."
        echo ""
        echo " Tail logs:"
        echo "   journalctl -u polaris.service -f"
        echo ""
        echo " Optional next steps:"
        echo "   - Install an ASTAP star catalog (V50 covers most rigs):"
        echo "     wget -O /tmp/v50.deb https://downloads.sourceforge.net/project/astap-program/star_databases/v50_star_database.deb"
        echo "     sudo dpkg -i /tmp/v50.deb"
        echo "   - Install GraXpert for AI post-processing:"
        echo "     see https://github.com/DanWBR/NINA.Polaris/blob/master/docs/user-guide/raspberry-pi-setup.md#44-graxpert-manual-install"
        echo "==============================================="
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
